# board/raspbery_pi/rootfs-overlay/etc/profile.d/golang.sh
export GOHOME=/usr/local/go
export PATH=$PATH:$GOHOME/bin
export PATH=$PATH:$(go env GOPATH)/bin
